# webhook/matchmaker-webhook.service
[Unit]
Description=MatchMaker GitHub Webhook Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/matchmaker/webhook
ExecStart=/usr/bin/python3 /opt/matchmaker/webhook/webhook_server.py
Restart=always
RestartSec=5
Environment=WEBHOOK_PORT=9000
Environment=DEPLOY_SCRIPT=/opt/matchmaker/deploy/auto-deploy.sh
# 安全：从 /opt/matchmaker/.env 读取 WEBHOOK_SECRET 和 BARK_KEY，避免在仓库中硬编码
EnvironmentFile=/opt/matchmaker/.env
Environment=PATH=/root/.nvm/versions/node/v24.15.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
